# supervisor_nginx.ini
[program:nginx]
command = /usr/bin/nginx -g 'pid /run/nginx.pid; error_log stderr;'
user = root
autostart = true
