# supervisor_php.ini
[program:php]
command = /usr/bin/php-fpm --nodaemonize --pid /run/php-fpm/php-fpm.pid
user = root
autostart = true
